# scepg.timer
[Unit]
Description=scepg timer

[Timer]
OnCalendar=*-*-* 00,06,12,18:00:00
Persistent=true

[Install]
WantedBy=multi-user.target
